# service/fossil/run
#!/bin/sh

exec /usr/bin/fossil server --repolist /fossil-data/fossils