# service/fossil/run
#!/bin/sh

if [ ! -e /fossil-data/fossils ]
then
   mkdir /fossil-data/fossils
fi
   
exec /usr/bin/fossil server --https --repolist /fossil-data/fossils
